# systemd/fpv-scan.service
[Unit]
Description=FPV HackRF scan service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/fpv-video-stream/agent/scan
Environment=SCAN_ID=scan-01
Environment=SCAN_SOURCE=live
# SDR backend: "bladerf" (default in config.py) drives a bladeRF 2.0 micro; "hackrf"
# falls back to hackrf_sweep/hackrf_transfer. See docs/bladerf-scanner-deploy.md.
Environment=SCAN_SDR=bladerf
Environment=SCAN_STATE_PATH=/run/fpv-scan/scan.json
# MQTT broker (scan data over WireGuard). Host is non-secret; the publish creds
# (MQTT_PUB_USER=pub, MQTT_PUB_PASS=...) are secret — put them in /etc/fpv-scan.env
# (not committed). The leading '-' lets the unit start even if that file is absent.
Environment=SCAN_MQTT_HOST=10.8.0.1
EnvironmentFile=-/etc/fpv-scan.env
# Analog-video frame extraction (reuses the live MQTT publisher -> fpv/<id>/video).
# Enabled by default; set FPV_VIDEO_ENABLED=0 to disable, FPV_EMIT_COOLDOWN_S to retune
# the per-channel throttle (default 10s). Full-res PNGs land in FPV_FRAMES_DIR.
Environment=FPV_FRAMES_DIR=/var/lib/fpv/frames
RuntimeDirectory=fpv-scan
StateDirectory=fpv
ExecStart=/opt/fpv-video-stream/agent/scan/.venv/bin/python -u main.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
